# ChargeLimiter/Package_rootless/DEBIAN/postinst
#!/bin/sh

APP_DIR="/var/jb/Applications/ChargeLimiter.app"
DAEMON_PLIST="/var/jb/Library/LaunchDaemons/com.chargelimiter.mod.plist"

chown -R root:wheel "$APP_DIR"
chmod +s "$APP_DIR/ChargeLimiterDaemon"
killall -9 ChargeLimiter ChargeLimiterDaemon 2>/dev/null

chown root:wheel "$DAEMON_PLIST"
launchctl bootout system "$DAEMON_PLIST" 2>/dev/null
launchctl bootstrap system "$DAEMON_PLIST" 2>/dev/null || launchctl load "$DAEMON_PLIST" 2>/dev/null

uicache -p "$APP_DIR" 2>/dev/null

exit 0
